# backup-tools/repl-image/files/crontab
5 0,8,16 * * * root test -x /usr/local/sbin/backup-replication.sh && /usr/local/sbin/backup-replication.sh > /proc/1/fd/1 2>&1
