# backup-tools/backup-image/files/crontab
*/5 * * * * root test -x /usr/local/sbin/db-backup.sh && /usr/local/sbin/db-backup.sh > /proc/1/fd/1 2>&1
